# systemd-sysusers.service.d/override.conf
# https://github.com/arkane-systems/genie/wiki/Systemd-units-known-to-be-problematic-under-WSL#systemd-sysusersservice
[Service]
LoadCredential=